# roles/llama_cpp/templates/llama-server.service.j2
[Unit]
Description=llama.cpp HTTP server ({{ llama_cpp_model_name }})
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ llama_cpp_server_user }}
Group={{ llama_cpp_server_user }}
Environment=LD_LIBRARY_PATH={{ llama_cpp_lib_dir }}
ExecStart={{ llama_cpp_install_dir }}/llama-server \
  --model {{ llama_cpp_model_dir }}/{{ llama_cpp_model_name }} \
  --alias {{ llama_cpp_server_model_alias }} \
  --host {{ llama_cpp_server_host }} \
  --port {{ llama_cpp_server_port }} \
  --ctx-size {{ llama_cpp_server_ctx_size }} \
  --threads {{ llama_cpp_server_threads }} \
  --parallel {{ llama_cpp_server_parallel }}
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
